# guacamole-client/rootfs/etc/services.d/tomcat/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Guacamole Client
# Runs Tomcat
# ==============================================================================

declare catalina_home="/usr/local/tomcat"

bashio::log.info "Starting Tomcat/Guacamole..."

cd "${catalina_home}" || bashio::exit.nok "Cannot change to Tomcat directory"

exec "${catalina_home}/bin/catalina.sh" run
